# market-dashboard/deploy/systemd/market-sync.service
[Unit]
Description=Market dashboard data sync
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/global_panel/market-dashboard
EnvironmentFile=/etc/market-dashboard.env
ExecStart=/root/global_panel/market-dashboard/.venv/bin/python /root/global_panel/market-dashboard/sync_market_data.py

