# market-dashboard/deploy/systemd/market-agent-report.timer
[Unit]
Description=Generate market dashboard AI research report every hour

[Timer]
OnBootSec=5min
OnActiveSec=5min
OnUnitActiveSec=60min
Unit=market-agent-report.service

[Install]
WantedBy=timers.target
