# market-dashboard/deploy/systemd/market-sync.timer
[Unit]
Description=Run market dashboard sync every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Unit=market-sync.service

[Install]
WantedBy=timers.target

